# system-maintenance/FreeBSD/feature/host.feature/files/_ROOT_/usr/local/etc/periodic/daily/00.app-update
#!/bin/sh
_LOGFILE=/var/log/walterjwhite/app-update.log app-update
