# system-maintenance/Linux/feature/host.feature/files/_ROOT_/etc/cron.daily/update-system
#!/bin/sh
[ -n "$_OPTN_SYSTEM_MAINTENANCE_DELAY_START" ] && {
    _WARN "delaying execution: $_OPTN_SYSTEM_MAINTENANCE_DELAY_START"
    sleep $_OPTN_SYSTEM_MAINTENANCE_DELAY_START
}
update-system
